# resources/djdash.conf
# dj dashbored
#

description	"dj dashboard"

start on runlevel [2345]
stop on runlevel [!2345]

setuid ken
setgid ken

respawn

exec java -jar /home/src/djdash-uber.jar
